# scripts/raspi-eye.service
# raspi-eye.service
# systemd service unit for raspi-eye smart camera application
#
# 安装方法:
#   sudo cp scripts/raspi-eye.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable raspi-eye.service
#   sudo systemctl start raspi-eye.service
#
# 查看日志:
#   journalctl -u raspi-eye.service -f

[Unit]
Description=RaspiEye Smart Camera Service
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=notify
ExecStart=/usr/local/bin/raspi-eye --config /etc/raspi-eye/config.toml
Environment=GST_PLUGIN_PATH=/usr/local/lib/raspi-eye/plugins
Restart=on-failure
RestartSec=5
WatchdogSec=30
TimeoutStartSec=60
TimeoutStopSec=35
StandardOutput=journal
StandardError=journal
SyslogIdentifier=raspi-eye
WorkingDirectory=/tmp
User=pi
Group=pi
SupplementaryGroups=video

# 安全加固
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/etc/raspi-eye

[Install]
WantedBy=multi-user.target
